# roles/environment/templates/profile.d/java.sh.j2
export JAVA_HOME={{ java_installdir }}/java
export PATH=$PATH:$JAVA_HOME/bin
