# roles/environment/templates/profile.d/maven.sh.j2
export M2_HOME={{ maven_installdir }}/apache-maven-3.0.5
export M2=$M2_HOME/bin 
export PATH=$M2:$PATH
